# board-bot.service
[Service]
WorkingDirectory=%h/board
EnvironmentFile=%h/board/.env
ExecStart=python3 discordbot.py
SyslogIdentifier=dancer.board
SyslogLevel=err
